# systemd/hcm-gis-cleanup.timer
[Unit]
Description=Daily cleanup of old HCM-GIS output files (2:00 AM Vietnam time)

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
